# ops/burner.service
[Unit]
Description=1000x Token Burner Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/1000x

# Main executable
ExecStart=/usr/bin/node /home/ubuntu/1000x/app/ts/burner_bot.js

# Restart policy
Restart=on-failure
RestartSec=30
StartLimitInterval=300
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=burner-bot

# Environment
EnvironmentFile=/home/ubuntu/1000x/ops/burner.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/home/ubuntu/1000x
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true

# Resource limits
MemoryMax=512M
MemoryHigh=256M
CPUQuota=50%

# Network restrictions (allow outbound for RPC calls)
RestrictAddressFamilies=AF_INET AF_INET6
IPAddressDeny=any
IPAddressAllow=localhost
IPAddressAllow=0.0.0.0/0

# Capability restrictions
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
